# deploy/alumni-backend.service
[Unit]
Description=Alumni Backend API
After=network.target mariadb.service

[Service]
Type=simple
User=nobody
WorkingDirectory=/app/backend
ExecStart=/app/backend/server
Restart=on-failure
RestartSec=5

Environment="DB_HOST=127.0.0.1"
Environment="DB_PORT=3306"
Environment="DB_USER=alumni"
Environment="DB_PASSWORD=password"
Environment="DB_NAME=alumni_db"
Environment="JWT_SECRET=change-me-in-production"
Environment="ALLOWED_ORIGIN=https://daeilfoundation.or.kr"
Environment="SITE_BASE_URL=https://daeilfoundation.or.kr"
Environment="KAKAO_CLIENT_ID="
Environment="KAKAO_CLIENT_SECRET="
Environment="KAKAO_REDIRECT_URI=https://daeilfoundation.or.kr/api/auth/kakao/callback"
Environment="UPLOAD_BASE_PATH=/var/www/uploads"
Environment="EASYPAY_GW_URL=gw.easypay.co.kr"
Environment="EASYPAY_RETURN_BASE_URL=https://www.daeilfoundation.or.kr"
Environment="PG_AUDIT_LOG_PATH=/var/logs/pg/pg-audit.log"
Environment="GOMEMLIMIT=80MiB"
Environment="GOMAXPROCS=1"

Environment="SMTP_HOST=smtp.naver.com"
Environment="SMTP_PORT=587"
Environment="SMTP_USER=ghkdwp018@naver.com"
Environment="SMTP_PASSWORD=9ZZSYSP3ZMPE"
Environment="SMTP_FROM=ghkdwp018@naver.com"

OOMScoreAdjust=0
MemoryLimit=120M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=alumni-backend

[Install]
WantedBy=multi-user.target
